# seo-workspace/config/daily-automation.cron.example
0 9 * * * cd '/Users/wangchao/Desktop/skill中心/skill-zhuangxiuseogeo' && 'python3' '.agents/skills/renovation-seo-geo/scripts/seo_geo_cli.py' 'daily-automation' '--pipeline' 'brief' '--website-root' '/Users/wangchao/Desktop/装修网站/zhuangxiuwangzhan-main' '--no-fetch-research-remote' '--research-search-provider' 'hybrid-rss' '--research-search-feeds-config' 'seo-workspace/config/research-search-feeds.example.yml' >> 'seo-workspace/reports/daily-automation.log' 2>&1
